# share/ocdc/examples/com.ocdc.poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ocdc.poll</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/ocdc</string>
        <string>poll</string>
        <string>--once</string>
    </array>
    
    <!-- Run every 5 minutes -->
    <key>StartInterval</key>
    <integer>300</integer>
    
    <!-- Also run at load -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>/tmp/ocdc-poll.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ocdc-poll.log</string>
    
    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/YOUR_USERNAME</string>
    </dict>
</dict>
</plist>
